# relay/deploy/yaver-relay.service
[Unit]
Description=Yaver Relay Server — P2P tunnel for mobile-to-desktop connectivity
Documentation=https://github.com/kivanccakmak/yaver.io
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# CONVEX_URL is REQUIRED for per-user password validation. Without it
# the relay can only validate against the shared RELAY_PASSWORD env
# var, and every web/mobile request that comes in with the user's
# per-user `__rp=...` (which doesn't match the shared password) gets
# rejected with 401 invalid relay password — even though Convex
# happily validates the same value via POST /relay/validate. The
# agent's QUIC register works because the agent sends the shared
# password, while the web SSE path sends the per-user one. Default
# points at the production Convex deployment; override via
# /etc/default/yaver-relay if a self-hosted relay needs a different
# Convex.
Environment=CONVEX_URL=https://perceptive-minnow-557.eu-west-1.convex.site
EnvironmentFile=-/etc/default/yaver-relay
ExecStart=/usr/local/bin/yaver-relay serve --quic-port=4433 --http-port=8443
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Logging — all output goes to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yaver-relay

# Performance
LimitNOFILE=65535

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
