# relay/deploy/yaver-relay.service
[Unit]
Description=Yaver Relay Server — P2P tunnel for mobile-to-desktop connectivity
Documentation=https://github.com/kivanccakmak/yaver.io
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# This unit ships with the OFFICIAL public.yaver.io deployment.
# YAVER_RELAY_OFFICIAL=1 tells the binary it's safe to fall back to
# Yaver's prod Convex backend if no CONVEX_URL is explicitly set.
# Self-hosted users running their own relay should NOT use this unit
# unchanged — clone it, drop YAVER_RELAY_OFFICIAL, and either set
# CONVEX_URL to their own Convex deploy or leave it empty for
# shared-password-only mode.
#
# CONVEX_URL is what actually drives per-user `__rp=` validation:
# the relay POSTs to {CONVEX_URL}/relay/validate when the inbound
# password doesn't match RELAY_PASSWORD. The agent's QUIC register
# works either way because the agent always sends the shared
# password; only the web/mobile per-user path needs Convex. Override
# CONVEX_URL via /etc/default/yaver-relay for staging-Convex deploys.
Environment=YAVER_RELAY_OFFICIAL=1
Environment=CONVEX_URL=https://perceptive-minnow-557.eu-west-1.convex.site
# EXPOSE_DOMAIN drives the auto-provisioned per-device subdomain that
# the relay assigns to every connected agent (e.g.
# <deviceId>.dev.yaver.io) and publishes back as `publicUrl` in the
# heartbeat. Routed via Cloudflare DNS (`*.dev.yaver.io` → A → this
# relay's IP, DNS-only / gray cloud) so the wildcard never enters the
# Worker that serves the apex `yaver.io` marketing site. Self-hosted
# operators override this in /etc/default/yaver-relay or leave it
# empty to disable the auto-subdomain feature entirely.
Environment=EXPOSE_DOMAIN=dev.yaver.io
EnvironmentFile=-/etc/default/yaver-relay
ExecStart=/usr/local/bin/yaver-relay serve --quic-port=4433 --http-port=8443 --expose-domain=${EXPOSE_DOMAIN}
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Logging — all output goes to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yaver-relay

# Performance
LimitNOFILE=65535

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
