# relay/deploy/yaver-relay-rotate.timer
[Unit]
Description=Monthly Yaver relay shared-password rotation

[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
