# pi-image/rootfs/etc/systemd/system/yaver-agent.service
[Unit]
Description=Yaver Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yaver
Group=yaver
Environment=HOME=/home/yaver
# /var/lib/yaver/workspaces is owned by yaver:yaver (set by the pi-image
# preseed in pi-image/scripts/setup-yaver-user.sh). Don't switch this back
# to /root or HOME=/root — repo paths leak the user's homedir name into
# Convex device rows, and the privacy test forbids it.
# --ghost enables the UI ghost so this appliance can drive a web-UI ERP (DIA /
# Tiger Wings / j-Platform / ERPNext) via its headless browser for migration +
# keeping the local ERP synced. Native-desktop ERPs (Logo 3 / Mikro / Netsis)
# are driven on a Windows Slave instead; this box commands it over the mesh.
ExecStart=/usr/local/bin/yaver serve --multi-user --ghost --port 18080 --work-dir /var/lib/yaver/workspaces --debug
Restart=always
RestartSec=5
WorkingDirectory=/var/lib/yaver
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
