# ci/remote/smoke/systemd/yaver-smoke-relay-password.service
[Unit]
Description=Yaver smoke check — relay password / fresh-install signup pipeline
Documentation=https://github.com/kivanccakmak/yaver.io/blob/main/ci/remote/smoke/relay-password.sh
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# /etc/yaver-smoke.env may override CONVEX_URL and SMOKE_TIMEOUT.
# It is created by ci/remote/smoke/install.sh and is NOT tracked in git.
EnvironmentFile=-/etc/yaver-smoke.env
ExecStart=/usr/local/libexec/yaver-smoke/relay-password.sh
# Keep the box polite to Convex: no retries, single attempt per tick.
TimeoutStartSec=90
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
# Hardening — this process only makes outbound HTTPS calls.
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
MemoryMax=128M

[Install]
WantedBy=multi-user.target
