# ci/remote/smoke/systemd/yaver-smoke-relay-password.timer
[Unit]
Description=Run Yaver relay-password smoke check every 15 minutes
Documentation=https://github.com/kivanccakmak/yaver.io/blob/main/ci/remote/smoke/relay-password.sh

[Timer]
# First run 2 min after boot so bootstrap has time to finish, then every
# 15 min. Small random delay so multiple smoke timers don't storm Convex.
OnBootSec=2min
OnUnitActiveSec=15min
RandomizedDelaySec=60
AccuracySec=30s
Persistent=true
Unit=yaver-smoke-relay-password.service

[Install]
WantedBy=timers.target
